# deploy/kriegspiel-simple-heuristics-bot.service
[Unit]
Description=Kriegspiel Simple Heuristics Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=fil
WorkingDirectory=/home/fil/dev/kriegspiel/bot-simple-heuristics
ExecStart=/home/fil/dev/kriegspiel/bot-simple-heuristics/.venv/bin/python /home/fil/dev/kriegspiel/bot-simple-heuristics/bot.py --poll-seconds 2
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
